# services/xset-on.service
[Unit]
Description=Xset led on

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/env xset led on

[Install]
WantedBy=default.target
